# rootfs/etc/services.d/redis-server/run
#!/usr/bin/with-contenv sh

exec redis-server /etc/redis/redis.conf --daemonize no
